# home/.config/autostart/polybar.desktop
[Desktop Entry]
Name=Polybar
Exec=polybar gnome
Terminal=false
Type=Application
StartupNotify=false
